# logid/99-logid-restart.rules
ACTION=="add", SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", RUN+="/usr/bin/systemctl restart logid"
